# .config/systemd/user/wallpaper.service
[Unit]
Description=Desktop wallpaper slideshow

[Service]
Type=oneshot
ExecStart=/bin/sh %h/bin/wallpaper
