# .config/systemd/user/wallpaper.timer
[Unit]
Description=Changes Wallpaper every hour
PartOf=graphical-session.target

[Timer]
OnBootSec=15m
OnUnitActiveSec=1h
