# .config/systemd/user/wallpaper.timer
[Unit]
Description=Changes Wallpaper every hour

[Timer]
OnCalendar=*-*-* *:00:00

[Install]
WantedBy=timers.target
