# .config/systemd/user/wm.target
[Unit]
Description=X server session for user
Requires=graphical-session.target
