# backups/backup.service
[Unit]
Description=Borg Backup

[Service]
ExecStart=/usr/bin/backup
