# backups/backup.timer
[Unit]
Description=Periodic Backup

[Timer]
OnCalendar=daily
Persistent=true
